# res/os/nf.service
[Unit]
Description=nf
[Service]
ExecStart=/usr/bin/weston
Environment="XDG_RUNTIME_DIR=/run/weston"
Type=idle
[Install]
WantedBy=multi-user.target
