# res/os/nf.service
[Unit]
Description=nf
[Service]
ExecStart=weston
Type=idle
[Install]
WantedBy=multi-user.target
