# res/os/nf.service
[Unit]
Description=nf
[Service]
ExecStart=/usr/bin/weston
Type=idle
[Install]
WantedBy=graphical.target
