# res/os/nth.service
[Unit]
Description=nth
[Service]
ExecStart=/usr/bin/python3 /root/os/site.py
Type=idle
[Install]
WantedBy=multi-user.target
