# res/ocuros.service
[Unit]
Description=ocurOS
[Service]
ExecStart=startx
Type=idle
[Install]
WantedBy=multi-user.target
