# src/config/bros.service
[Unit]
Description=bros
[Service]
ExecStart=/usr/bin/bros -n
Type=idle
[Install]
WantedBy=graphical.target