# src/sh/brs.service
[Unit]
Description=brs
[Service]
ExecStart=/usr/bin/python3 /root/nfnth/src/py/site.py
Type=idle
[Install]
WantedBy=multi-user.target
